# deploy/openrc/pulse-node
#!/sbin/openrc-run

name="pulse-node"
description="Pulse Node"
command="/usr/local/bin/pulse-node"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
directory="/var/lib/pulse"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}

start_pre() {
	if [ -f /etc/pulse/pulse-node.env ]; then
		set -a
		. /etc/pulse/pulse-node.env
		set +a
	fi
}
